# server/systemd/igitt.socket
[Unit]
Description=IGITT server socket

[Socket]
ListenStream=[::1]:15177
Accept=false
Service=igitt.service

[Install]
WantedBy=igitt.service
